# scripts/pru/pickupwinder-pru.service
[Unit]
Description=PickupWinder PRU firmware loader
Wants=systemd-udev-settle.service
After=local-fs.target systemd-modules-load.service systemd-udev-settle.service
ConditionPathExists=/lib/firmware/am335x-pru0-fw
ConditionPathExists=/lib/firmware/am335x-pru1-fw

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/load_pru.sh --start-only
ExecStop=/usr/local/bin/load_pru.sh --stop-only
TimeoutStartSec=30
TimeoutStopSec=15
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
